# .docker/os2web/crontab
* * * * * drush cron --root=/opt/drupal >> /var/log/cron.log 2>&1
# Mandatory blank line
